# nws-alert-poller.service
[Unit]
Description=NWS Alert Poller
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
# IMPORTANT: change you@domain.com to a real email address (NWS policy)
Environment=NWS_USER_AGENT=FreePBX-NWS-Alert/1.0 (contact: you@domain.com)
# Optional: adjust pre-wait before playback (seconds)
Environment=NWS_PREWAIT_SEC=2
# Optional: delay between sequential alerts (seconds) - prevents overlapping playback
Environment=NWS_ALERT_DELAY_SEC=30
# Caller ID shown on alerted phones (name and number)
Environment=NWS_CID_NAME=System Alert
Environment=NWS_CID_NUM=0000
ExecStart=/usr/bin/python3 /usr/local/bin/nws_alert_poller.py
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
